# Docker/rootfs/etc/services.d/NetDaemonAddOnApp/run
#!/usr/bin/with-contenv /usr/bin/bashio

# Optional MQTT configuration
declare Mqtt__Host
declare Mqtt__Port
declare Mqtt__Username
declare Mqtt__Password

# Set configuration values to environment variables
NetDaemon__ApplicationAssembly=$(bashio::config 'app_assembly')
Logging__Loglevel__Default=$(bashio::config 'log_level') 
NetDaemon__ApplicationConfigurationFolder=$(bashio::config 'app_config_folder')

if bashio::services.available "mqtt"; then
    Mqtt__Host=$(bashio::services 'mqtt' 'host')
    Mqtt__Port=$(bashio::services 'mqtt' 'port')
    Mqtt__Username=$(bashio::services 'mqtt' 'username')
    Mqtt__Password=$(bashio::services 'mqtt' 'password')
else
    bashio::log.warning \
        "No MQTT add-on installed, MQTT features will not work."
fi

export \
    HomeAssistant__Host="supervisor" \
    HomeAssistant__WebsocketPath="core/websocket" \
    HomeAssistant__Port=80 \
    HomeAssistant__Ssl=false \
    HomeAssistant__Token="${SUPERVISOR_TOKEN}" \
    Mqtt__Host \
    Mqtt__Port \
    Mqtt__Username \
    Mqtt__Password \
    NetDaemon__ApplicationAssembly \
    Logging__Loglevel__Default \
    NetDaemon__ApplicationConfigurationFolder

declare daemondir="/daemon"

if ! bashio::fs.file_exists "/data/options.json"; then
    bashio::exit.nok "This container only supports running as Home Assistant add-on!"
fi

if ! bashio::config.has_value "app_assembly"; then
    bashio::log.info "Starting NetDaemon V3 Runtime ..."
    cd "${daemondir}" || bashio::exit.nok "${daemondir} does not exist"
    exec dotnet NetDaemon.Host.Default.dll
else
    # We have provided an application assembly setting
    if [[ "${NetDaemon__ApplicationAssembly}" == *".csproj" ]]; then
        bashio::exit.nok "csproj deployments are not supported in v3, use compiled option instead!"
    fi

    if [[ "${NetDaemon__ApplicationAssembly}" != *".dll" ]]; then
        bashio::exit.nok "Assembly needs to point to a .dll file!"
    fi

    # make path relative to data folder (/config/netdaemon if addon) 
    # if the path is a relative path
    if [[ "${NetDaemon__ApplicationAssembly}" != "/"* ]]; then
        export NetDaemon__ApplicationAssembly="/config/netdaemon3/${NetDaemon__ApplicationAssembly}"
    fi 

    # The provided application source is ether a project or pre-compiled .Net application
    if ! bashio::fs.file_exists "${NetDaemon__ApplicationAssembly}"; then
        bashio::exit.nok "The assembly ${NetDaemon__ApplicationAssembly} cannot be found. Please check the settings."
    fi
    
    bashio::log.info "Starting NetDaemon V3 pre-built Runtime using assembly ${NetDaemon__ApplicationAssembly}..."
    cd "$(dirname "${NetDaemon__ApplicationAssembly}")" ||  bashio::exit.nok "Could not change directory to run project"

    exec dotnet "$(basename "${NetDaemon__ApplicationAssembly}")"
fi
