# Docker/rootfs/etc/services.d/NetDaemonAdmin/run
#!/usr/bin/with-contenv bash
echo "Starting NetDaemon Admin"

cd /admin || exit 1

exec node ./host/admin.js