# Docker/rootfs/etc/services.d/NetDaemon/run
#!/bin/bash

# ==============================================================================
# NetDaemon Service
# Starts NetDaemon
# ==============================================================================
cd $HASS_RUN_PROJECT_FOLDER || echo -e "\033[31mCould not change directory to run project\033[0m" >&2

echo -e "\033[32mStarting NetDaemon...\033[0m" >&2

dotnet run
